# pdf-downloader.service
[Unit]
Description=PDF Paper Downloader API Service
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/opt/pdf-downloader
Environment="PATH=/opt/pdf-downloader/venv/bin"
ExecStart=/opt/pdf-downloader/venv/bin/gunicorn --workers 4 --bind 0.0.0.0:5000 wsgi:app
Restart=always

[Install]
WantedBy=multi-user.target